# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for tdd-lists-staging"

start on net-device-up
stop on shutdown

respawn

setuid ubuntu
chdir /home/ubuntu/sites/tdd-lists-staging/source/tddsuperlists

exec ../../virtualenv/bin/gunicorn --bind unix:/tmp/tdd-lists-staging.socket superlists.wsgi:application
